# novaic-quic-service/deploy/novaic-quic-service.service
[Unit]
Description=NovAIC QUIC Service (STUN + Relay)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/novaic/novaic-quic-service

# 环境变量（按需修改）
Environment=GATEWAY_URL=https://api.gradievo.com
Environment=RELAY_PORT=443
Environment=STUN_PORT=3478
Environment=RELAY_TLS_CERT_PATH=/etc/letsencrypt/live/relay.gradievo.com/fullchain.pem
Environment=RELAY_TLS_KEY_PATH=/etc/letsencrypt/live/relay.gradievo.com/privkey.pem
Environment=RUST_LOG=info
# RUST_LOG=debug 可查看连接尝试（排查 VNC 连接问题时临时启用）

ExecStart=/opt/novaic/novaic-quic-service/novaic-quic-service
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
